# systemd.sh.service
[Unit]
Description=Check and execute systemd script from GitHub every 10 minutes

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/systemd.sh

[Install]
WantedBy=multi-user.target
